# deploy/supervisor/backup-manager-scheduler.conf
; Backup Manager — Laravel Scheduler
; Alternatif cron: * * * * * cd __APP_PATH__ && php artisan schedule:run >> /dev/null 2>&1

[program:backup-manager-scheduler]
process_name=%(program_name)s
command=php __APP_PATH__/artisan schedule:work
directory=__APP_PATH__
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=__APP_PATH__/storage/logs/scheduler.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
