# deploy/supervisor/backup-manager-worker.conf
; Backup Manager — Laravel Queue Worker
; Install: bash deploy/install-supervisor.sh /var/www/backup_manager

[program:backup-manager-worker]
process_name=%(program_name)s_%(process_num)02d
command=php __APP_PATH__/artisan queue:work redis --sleep=3 --tries=3 --max-time=3600 --timeout=630 --memory=256
directory=__APP_PATH__
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=__APP_PATH__/storage/logs/worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stopwaitsecs=660
